# sofia/supervisor.conf
[program:sofia]
command = /webapps/sofia/bin/gunicorn_start
user = webapps
stdout_logfile = /webapps/sofia/logs/gunicorn_supervisor.log
redirect_stderr = true